# root/ram/etc/inittab
id:2:initdefault:

si:S:sysinit:/sbin/python -m boot
